# services/spotify-display.service
[Unit]
Description=Spotify Pi Display — Flask Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/spotify-pi
ExecStart=/home/pi/spotify-pi/venv/bin/python server.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
