# services/spotify-display.service
[Unit]
Description=Spotify Pi Display — Flask Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=admin
WorkingDirectory=/home/admin/circle-pi-display
ExecStart=/home/admin/circle-pi-display/venv/bin/python server.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
