# services/spotify-display.service
[Unit]
Description=Spotify Pi Display — Flask Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__DEPLOY_USER__
WorkingDirectory=__DEPLOY_DIR__
ExecStart=__DEPLOY_DIR__/venv/bin/python server.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
